# docker/prod/supervisord_cron.conf
[supervisord]
nodaemon=false

[program:cron]
command=cron -f
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
